# scripts/services/mysterium-update.timer
[Unit]
Description=Timer runs mysterium-update.service Daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target